# desktop/.config/rc/init.d/sfwbar
#!/sbin/openrc-run

description="Sway Floating Window Bar"

command="/usr/local/bin/sfwbar-git"

supervisor=supervise-daemon
error_log="${HOME}/.local/state/log/${RC_SVCNAME}.log"

depend() {
	need pipewire
}
