# desktop/.config/rc/init.d/wapanel
#!/sbin/openrc-run

description="Simple panel/status bar/task bar."

command="/usr/local/bin/wapanel"

supervisor=supervise-daemon
error_log="${HOME}/.local/state/log/${RC_SVCNAME}.log"

depend() {
	need pipewire
	need dbus
}
