# desktop/.config/rc/init.d/wapanel
#!/sbin/openrc-run

description="Simple panel/status bar/task bar."

command="/usr/local/bin/wapanel"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

depend() {
	need pipewire
	need dbus
}
