# desktop/.config/rc/init.d/waybar
#!/sbin/openrc-run

description="Highly customizable Wayland bar."

command="/usr/bin/waybar"

supervisor=supervise-daemon
error_log="${HOME}/.local/state/log/${RC_SVCNAME}.log"

depend() {
	before *
	need pipewire
}
